# ops/systemd/oper-worker-digest.service
[Unit]
Description=OnePercent Saved-Search Digest Worker
After=oper-postgres.service
Requires=oper-postgres.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/onepercent
EnvironmentFile=/etc/oper.env
ExecStart=/usr/bin/node apps/worker/dist/digest.js
Restart=on-failure
RestartSec=5
MemoryMax=192M

[Install]
WantedBy=multi-user.target
