# ops/systemd/oper-worker-alerts.service
[Unit]
Description=OnePercent Pro Deal Flow Alert Worker
After=oper-postgres.service
Requires=oper-postgres.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/onepercent
EnvironmentFile=/etc/oper.env
EnvironmentFile=-/etc/oper-role-worker.env
# alerts-main does NOT import @oper/query-lang, so it runs through plain node
# (dist is emitted by `pnpm --filter @oper/worker build` -> tsc). Do NOT use tsx
# here — that runtime is reserved for digest/watchlist which pull in query-lang.
# tsx is a prod dep of @oper/worker. WorkingDirectory is /opt/onepercent.
ExecStart=/usr/bin/node apps/worker/dist/alerts-main.js
Restart=on-failure
RestartSec=5
# Fanout holds the tick's candidates + watchlists in memory; 192M could be
# approached at large matched-row volume (#42). 384M gives headroom while the
# ALERT_WATCHLIST_BATCH cap (#41) keeps the watchlist set bounded per tick.
MemoryHigh=288M
MemoryMax=384M

[Install]
WantedBy=multi-user.target
