# ops/systemd/oper-scraper.service
[Unit]
Description=OnePercent Scraper Service (FastAPI geocoding)
After=oper-postgres.service
Requires=oper-postgres.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/onepercent
EnvironmentFile=/etc/oper.env
EnvironmentFile=-/etc/oper-role-worker.env
Environment="PYTHONPATH=/opt/onepercent/services/scraper_service:/opt/onepercent/services"
ExecStart=/opt/onepercent/services/ml/.venv/bin/uvicorn \
    services.scraper_service.main:app \
    --host 127.0.0.1 \
    --port 8001
Restart=on-failure
RestartSec=5
# homeharvest holds the full result DataFrame in memory and, with
# extra_property_data, fan-out-fetches per-property details (20 workers). Dense
# metros + the combined for_sale+for_rent pass push peak RSS well past 1G and
# were OOM-killed (2026-07-14). MemoryHigh throttles (reclaim) before the hard
# MemoryMax kill; sized for a ~4G box.
MemoryHigh=1536M
MemoryMax=2G

[Install]
WantedBy=multi-user.target
