# ops/systemd/oper-app.service
[Unit]
Description=OnePercent Next.js Frontend (app, port 3001)
After=oper-postgres.service oper-redis.service
Requires=oper-postgres.service oper-redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/onepercent/apps/one/.next/standalone
EnvironmentFile=/etc/oper.env
EnvironmentFile=-/etc/oper-role-app.env
Environment=PORT=3001
ExecStart=/usr/bin/node apps/one/server.js
Restart=on-failure
RestartSec=5
MemoryHigh=1536M
MemoryMax=2G

[Install]
WantedBy=multi-user.target
