# ops/systemd/oper-worker.service
[Unit]
Description=OnePercent Crawl Worker
After=oper-postgres.service oper-redis.service
Requires=oper-postgres.service oper-redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/onepercent
EnvironmentFile=/etc/oper.env
EnvironmentFile=-/etc/oper-role-worker.env
ExecStart=/usr/bin/node apps/worker/dist/crawl.js
Restart=on-failure
RestartSec=5
MemoryMax=512M

[Install]
WantedBy=multi-user.target
